# server/deploy/tedi-remote.service
[Unit]
Description=TEDI remote-access relay (agent <-> browser bridge)
After=network.target

[Service]
Type=simple
# Your Linux account that owns /var/www/html/tedi-remote.
User=<your-user>
WorkingDirectory=/var/www/html/tedi-remote/server
EnvironmentFile=/var/www/html/tedi-remote/server/.env
# Point ExecStart at your node binary (an nvm path is shown). Update it if you
# upgrade node, or symlink a stable path.
ExecStart=/home/<your-user>/.nvm/versions/node/<version>/bin/node /var/www/html/tedi-remote/server/server.js
Restart=always
RestartSec=2
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
