# server/deploy/tedi-remote.service
[Unit]
Description=TEDI remote-access relay (agent <-> browser bridge)
After=network.target

[Service]
Type=simple
# Your Linux account that owns /var/www/html/tedi-remote.
User=<your-user>
WorkingDirectory=/var/www/html/tedi-remote/server
EnvironmentFile=/var/www/html/tedi-remote/server/.env
# Point ExecStart at your node binary (an nvm path is shown). Update it if you
# upgrade node, or symlink a stable path.
ExecStart=/home/<your-user>/.nvm/versions/node/<version>/bin/node /var/www/html/tedi-remote/server/server.js
Restart=always
RestartSec=2

# Hardening — the relay is the only public-facing service, writes nothing to
# disk, needs no privileges, and only talks TCP on localhost. (MemoryDenyWrite
# Execute / SystemCallFilter are intentionally omitted: they break V8's JIT.)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=
# Resource caps also blunt a connection flood at the OS level.
MemoryMax=256M
TasksMax=128
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
